# docker/cron/crontab
*/2 * * * * php /var/www/scripts/file.php >> /var/log/cron.log 2>&1